# dnsmasq-addon/rootfs/etc/services.d/dnsmasq/run
#!/command/with-contenv bashio
set -euo pipefail
bashio::log.info "Avvio dnsmasq..."
exec dnsmasq -k --clear-on-reload
